# recipes-extended/logrotate/logrotate/logrotate.crontab
*/5 * * * * root /usr/sbin/logrotate /etc/logrotate.conf
